# files/wm/common/greetd/usr/lib/systemd/system/greetd.service.d/10-custom.conf
[Unit]
Conflicts=getty@tty1.service

[Service]
Type=idle
