# files/wm/common/wayvnc/usr/lib/systemd/user/wayvnc.socket
[Unit]
Description=WayVNC socket

[Socket]
ListenStream=0.0.0.0:5900

[Install]
WantedBy=sockets.target
