# files/wm/common/flatpak/usr/lib/systemd/system/flatpak-update.timer
[Unit]
Description=Daily Flatpak updates
Requires=flatpak-update.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target