# files/wm/common/flatpak/usr/lib/systemd/system/flatpak-update.service
[Unit]
Description=Update Flatpak applications
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak-update

[Install]
WantedBy=multi-user.target