# files/base/system/usr/lib/systemd/system/bootc-update-check.timer
[Unit]
Description=Check for bootc updates periodically

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
